# schedule/root/etc/supervisor/conf.d/publisher.conf
[program:schedule]
command=python /app/main.py
stderr_logfile=/dev/stderr
stdout_logfile=/dev/stdout
startsecs=0
