# publisher/root/etc/supervisor/conf.d/privoxy.conf
[program:privoxy]
command=privoxy --no-daemon /etc/privoxy/config
startsecs=0
