# publisher/root/etc/supervisor/conf.d/tor.conf
[program:tor]
command=/usr/bin/tor
startsecs=0
